# supervisor.conf
[program:elon-bot-api]
command=/path/to/elon-bot/venv/bin/python /path/to/elon-bot/api.py
directory=/path/to/elon-bot
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/path/to/elon-bot/logs/api.log
stderr_logfile=/path/to/elon-bot/logs/api-error.log

[program:elon-bot-bot]
command=/path/to/elon-bot/venv/bin/python /path/to/elon-bot/bot.py
directory=/path/to/elon-bot
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/path/to/elon-bot/logs/bot.log
stderr_logfile=/path/to/elon-bot/logs/bot-error.log
